# services/jackett/run
#!/bin/bash

/sbin/setuser abc mono /app/Jackett/JackettConsole.exe
