# root/etc/services.d/jackett/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc mono /app/Jackett/JackettConsole.exe "${RUN_OPTS}"
