# root/etc/services.d/jackett/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /app/Jackett/jackett_launcher.sh "${RUN_OPTS}"
